# files/iostat/crontab
# Collect data iops for zabbix
* * * * * root /usr/bin/timeout 60 /usr/libexec/zabbix-extensions/scripts/iostat-collect.sh /tmp/iostat.out 59 > /dev/null
