# files/iostat/crontab
# Collect data iops for zabbix - just execute "crontab -e" and write this line:
* * * * * DATA= $(/usr/bin/iostat -yxd -o JSON 59 1) ; echo $DATA >/tmp/iostat-cron.out
