# files/iostat/crontab
# Collect data iops for zabbix - just execute "crontab -e" and write this line:
* * * * * DATA=$(S_TIME_FORMAT=ISO /usr/bin/iostat -yxdt -o JSON 59 1) ; echo "$DATA" >/tmp/iostat-cron.out
