# files/iostat/crontab
# Collect data iops for zabbix - just execute "crontab -e" and write this line:
* * * * * /usr/libexec/zabbix-extensions/scripts/iostat-collect.sh /tmp/iostat-cron.out 59 >/dev/null 2>&1
