# files/iostat/crontab
# Collect data iops for zabbix - just execute "crontab -e" and write this line:
* * * * * /usr/bin/iostat -ydt -o JSON 59 1 >/tmp/iostat-cron.out
