# files/iostat/crontab
# Collect data iops for zabbix
* * * * * root /usr/libexec/zabbix-extensions/scripts/iostat-collect.sh /tmp/iostat.out 60
