# files/hwraid-smartarray/crontab
# Описание: сбор информации о HP SmartArray и отправка на сервер мониторинга
*/10    *       *       *       *       /etc/zabbix/scripts/hp-raid-data-processor.sh
