# files/linux/zabbix.cron
# Author:	Lesovsky A.V.
# Description:	Extended zabbix check
0 */2 * * *	root /etc/zabbix/scripts/check-open-descriptors.sh
