# files/linux/zabbix.cron
# Author:	Lesovsky A.V.
# Description:	Extended zabbix check
0 */2 * * *	root /usr/libexec/zabbix-extensions/scripts/check-open-descriptors.sh
