# files/linux/zabbix.cron
# Описание:     дополнительные проверки для Zabbix
# Автор:        Лесовский А.В.
# Добавил:      05-07-2013
0 */2 * * *	root /etc/zabbix/scripts/check-open-descriptors.sh
